# deployment/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:street]
command=npm run server

[program:cms]
command=npm run server-cms

[program:puppeteer]
command=npm run server-puppeteer
